# event-signup.service
[Unit]
Description=Event Signup Application
After=network.target

[Service]
Type=simple
User=deployer
Group=deployer
WorkingDirectory=/home/deployer/event-signup
ExecStart=/home/deployer/event-signup/event-signup-app
Restart=on-failure
RestartSec=5

# Environment
EnvironmentFile=/home/deployer/event-signup/.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/deployer/event-signup
PrivateTmp=true

[Install]
WantedBy=multi-user.target
